# services/infy_db_service/infy_db_service.service
[Unit]
Description=Infosys DB Service (Dev)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/projadmin/workarea/docwbsln/services/infy_db_service/src
ExecStart=/usr/bin/bash -c "source /home/projadmin/workarea/docwbsln/services/infy_db_service/.venv/bin/activate; export $(grep -v '^#' ../.env | xargs -d '\n'); python main.py"
User=projadmin
Group=projadmin
Restart=always
RestartSec=9
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=InfyDBService

[Install]
WantedBy=multi-user.target